# tests/syntax_test_system.crontab
# SYNTAX TEST "Packages/Crontab/System Crontab.sublime-syntax"

# crontab that controls running this proc
# ^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^ comment.line.number-sign
# TODO - move running to systemd

#login as root
#create job with preferred editor (e.g. vim)
MAILTO=root
#^^^^^ variable.other.readwrite.crontab
#     ^ keyword.operator.assignment.crontab
#      ^^^^ meta.string.crontab string.unquoted.crontab

* * * * * root touch /tmp/file
#^^^^^^^^ meta.string.cron-expression.crontab
# ^ meta.sequence.cron.hour.crontab constant.other.wildcard.asterisk.crontab
#   ^ meta.sequence.cron.day-of-month.crontab constant.other.wildcard.asterisk.crontab
#     ^ meta.sequence.cron.month.crontab constant.other.wildcard.asterisk.crontab
#       ^ meta.sequence.cron.day-of-week.crontab constant.other.wildcard.asterisk.crontab
#         ^^^^ constant.language.crontab
#              ^^^^^^^^^^^^^^^ source source

# m h dom mon dow user      command
*   *  *   *   *  someuser  echo 'foo'
#^^^^^^^^^^^^^^^ meta.string.cron-expression.crontab
#   ^ meta.sequence.cron.hour.crontab constant.other.wildcard.asterisk.crontab
#      ^ meta.sequence.cron.day-of-month.crontab constant.other.wildcard.asterisk.crontab
#          ^ meta.sequence.cron.month.crontab constant.other.wildcard.asterisk.crontab
#              ^ meta.sequence.cron.day-of-week.crontab constant.other.wildcard.asterisk.crontab
#                 ^^^^^^^^ entity.name.user.crontab
#                           ^^^^^^^^^^ source source
