# tests/syntax_test_system.crontab
# SYNTAX TEST "Packages/Crontab/System Crontab.sublime-syntax"

# crontab that controls running this proc
# ^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^^ comment.line.number-sign
# TODO - move running to systemd

#login as root
#create job with preferred editor (e.g. vim)
MAILTO=root
#^^^^^ variable.other.readwrite.crontab
#     ^ keyword.operator.assignment.crontab
#      ^^^^ meta.string.crontab string.unquoted.crontab

* * * * * root touch /tmp/file
#^^^^^^^^ meta.string.cron-expression.crontab
# ^ meta.sequence.cron.hour.crontab variable.language.wildcard.crontab
#   ^ meta.sequence.cron.day-of-month.crontab variable.language.wildcard.crontab
#     ^ meta.sequence.cron.month.crontab variable.language.wildcard.crontab
#       ^ meta.sequence.cron.day-of-week.crontab variable.language.wildcard.crontab
#         ^^^^ constant.language.crontab
#              ^^^^^^^^^^^^^^^ source.shell.bash.crontab
#              ^^^^^ meta.function-call.identifier.shell meta.command.shell variable.function.shell
#                   ^^^^^^^^^^ meta.function-call.arguments.shell
#                    ^^^^^^^^^ meta.string.glob.shell string.unquoted.shell

# m h dom mon dow user      command
*   *  *   *   *  someuser  echo 'foo'
#^^^^^^^^^^^^^^^ meta.string.cron-expression.crontab
#   ^ meta.sequence.cron.hour.crontab variable.language.wildcard.crontab
#      ^ meta.sequence.cron.day-of-month.crontab variable.language.wildcard.crontab
#          ^ meta.sequence.cron.month.crontab variable.language.wildcard.crontab
#              ^ meta.sequence.cron.day-of-week.crontab variable.language.wildcard.crontab
#                 ^^^^^^^^ entity.name.user.crontab
#                           ^^^^^^^^^^ source.shell.bash.crontab
#                           ^^^^ meta.function-call.identifier.shell support.function.shell
#                               ^^^^^^ meta.function-call.arguments.shell
#                                ^^^^^ meta.string.glob.shell string.quoted.single.shell
#                                ^ punctuation.definition.string.begin.shell
#                                    ^ punctuation.definition.string.end.shell
